# helloworld.service
[Unit]
Description=Say Hello
[Service]
ExecStart=/usr/bin/logger -i Hello World